# contrib/docker/crontab
*/10 * * * * root cd /src/yggdrasil-map/web/ && python updateGraph.py
